# overlay/system/etc/systemd/user/ayatana-indicator-sound.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/pulseaudio -k
