# overlay/system/etc/systemd/system/bluebinder.service.d/override.conf
[Unit]
After=device-hacks.service

[Service]
ExecStartPost=/sbin/rfkill unblock all
