# overlay/system/etc/systemd/system/mount-android-partitions.service.d/override.conf
[Service]
ExecStartPost=/usr/bin/mount -o bind /mnt/vendor/persist /android/persist
